# configs/auto-update-panel.timer
[Unit]
Description=Run Hydrodactyl Panel Auto-Updater periodically

[Timer]
OnBootSec=5min
OnUnitActiveSec=1h
AccuracySec=1min

[Install]
WantedBy=timers.target
